# supervisord.conf
[supervisord]
nodaemon=true
pidfile=/run/supervisord.pid
logfile=/var/log/supervisord.log
childlogdir=/var/log/supervisor

[program:sshd]
command=/bin/sshd -D
priority=1
autorestart=true
